# ops/systemd/emfi-view.service.example
[Unit]
Description=EMFI View static dashboard
After=network-online.target emfi-control.service emfi-develop.service
Wants=network-online.target

[Service]
Type=simple
User=lab
Group=lab
WorkingDirectory=/home/lab/EM-Architecture/view
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/npm run preview -- --host 127.0.0.1 --port 8003
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
